# mqttcmdsrv.conf
[program:mqttcmdsrv]
command = /usr/bin/mqttcmdsrv 172.29.1.1 1883 "domoticz/out" /etc/mqttcmdtrig.json
priority = 110
